# docker/supervisor.conf
[supervisord]
nodaemon=true

[program:sso]
command=python /app/backend/server.py %(ENV_SSO_HOSTNAME) /app/plex/Preferences.xml

[program:nginx]
command=nginx -c /etc/nginx/nginx.conf -g 'daemon off;'